# filesystem/fs/etc/inittab
#
#	$Id: inittab,v 1.2 2006/04/27 21:40:04 jwessel Exp $
# Busybox's init uses a different syntax then the normal init.
# The "-" in -/bin/sh allows job control on /dev/console.
#
::sysinit:/etc/init.d/rcS
::askfirst:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
